# scripts/com.echolog.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.echolog.daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>--import</string>
    <string>tsx</string>
    <string>src/server/app.ts</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/example/code/0toy/self_record</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/echolog.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/echolog.stderr.log</string>
</dict>
</plist>
